# config/.config/systemd/user/battery-daemon.service
[Unit]
Description=Battery Monitor Daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=/bin/bash %h/.local/bin/battery-daemon.sh
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
